# certbot/crontab
0 0,12 * * * /usr/local/bin/renew_certs.sh